# docker/cicd/billrun.conf
[program:php-fpm]
command=/usr/local/sbin/php-fpm --nodaemonize --fpm-config /usr/local/etc/php-fpm.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=%(ENV_BILLRUN_WEB_ENABLED)s
directory=/billrun
autorestart=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;' -c /etc/nginx/nginx.conf
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=%(ENV_BILLRUN_WEB_ENABLED)s
autorestart=true
redirect_stderr=false

[program:billrun-worker]
command=php /billrun/public/index.php --env container --worker
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
autostart=%(ENV_BILLRUN_WORKER_ENABLED)s
directory=/billrun
autorestart=true
